# techstop/conf/command.conf
[program:command]
command         = /var/www/techstop/env.sh
