# techstop/conf/command.conf
[program:command]
command         = /django/env.sh
